# JangB-main/JangB-main/jangb.service
# JangB 장비 대여 앱 - systemd 서비스
# 설치: sudo cp ~/JangB-main/jangb.service /etc/systemd/system/
#       sudo systemctl daemon-reload && sudo systemctl enable jangb && sudo systemctl start jangb

[Unit]
Description=JangB Equipment Rental Node Server
After=network.target

[Service]
Type=simple
User=pje070163
WorkingDirectory=/home/pje070163/JangB-main
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=3
StandardOutput=append:/home/pje070163/JangB-main/jangb-node.log
StandardError=append:/home/pje070163/JangB-main/jangb-node.log

[Install]
WantedBy=multi-user.target
